# datawall-restart.service
[Unit]
Description=Datawall Auto-Restart Service
After=graphical.target network-online.target
Wants=graphical.target network-online.target

[Service]
Type=simple
User=pi
Environment="DISPLAY=:0"
Environment="XAUTHORITY=/home/pi/.Xauthority"
ExecStart=/opt/datawall/datawall-restart.sh
Restart=always
RestartSec=10

[Install]
WantedBy=graphical.target