# datawall-restart.service
[Unit]
Description=Datawall Auto-Restart Service
After=graphical.target network-online.target
Wants=graphical.target network-online.target

[Service]
Type=simple
User=admin
Environment="DISPLAY=:0"
Environment="XAUTHORITY=/home/admin/.Xauthority"
Environment="XDG_RUNTIME_DIR=/run/user/1000"
ExecStart=/opt/datawall/datawall-restart.sh
Restart=always
RestartSec=10

[Install]
WantedBy=graphical.target